# supervisor-slapd.conf
[program:slapd]
priority=20
directory=/tmp
command=/usr/bin/slapd-start
user=root
autostart=true
autorestart=true
stdout_logfile=syslog
stderr_logfile=syslog

